# tailscan.service
[Unit]
Description=TailScan — Tailscale Network Scanner
After=network.target tailscaled.service
Wants=tailscaled.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/tailscan
EnvironmentFile=/opt/tailscan/.env
ExecStart=/opt/tailscan/venv/bin/python /opt/tailscan/main.py
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal

# Allow nmap raw sockets (needs root or CAP_NET_RAW)
# If running as non-root, add: AmbientCapabilities=CAP_NET_RAW CAP_NET_ADMIN

[Install]
WantedBy=multi-user.target
