# tailscan.service
[Unit]
Description=TailScan — Tailscale LAN Scanner
After=network.target tailscaled.service
Wants=tailscaled.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/tailscan
ExecStart=/usr/bin/python3 /opt/tailscan/main.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
